# roles/dav/files/vdirsyncer-sync.service
[Unit]
Description=Sync Vdirsyncer

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
