# roles/taskwarrior/files/bugwarrior.service
[Unit]
Description=bugwarrior-pull

[Service]
Type=oneshot
ExecStart=/usr/bin/bugwarrior-pull
